# versioncheck/crontab
*/5 * * * * root /home/steam/ark/versioncheck/versioncheck.sh
# An empty line is required at the end of this file for a valid cron file.